# server/service/temp-server-mon.timer
[Unit]
Description=Run temp-server-mon.service every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
