# server/smonitor.timer
[Unit]
Description=Run smonitor.service every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
